# CMake/config/xrmd.service.in
[Unit]
Description=Xilinx Resource Manager Daemon
Wants=network-online.target
After=network-online.target
ConditionDirectoryNotEmpty=/dev/xfpga

[Service]
Type=simple
LimitNOFILE=1000000
Environment="LD_LIBRARY_PATH=/opt/xilinx/xrt/lib"
Restart=always
RestartSec=1
ExecStart=@CMAKE_INSTALL_PREFIX@/xrm/bin/xrmd

[Install]
WantedBy=multi-user.target
